# docker/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
pidfile=/tmp/supervisord.pid
childlogdir=/var/log/supervisor

; ---------------------------------------------------------------------------
; Embedded Redis — skipped when REDIS_EMBEDDED=false
; ---------------------------------------------------------------------------
[program:redis]
command=redis-server
    --save 60 1
    --dir /data/redis
    --bind 127.0.0.1
    --maxmemory 256mb
    --maxmemory-policy allkeys-lru
    --loglevel warning
autostart=%(ENV_REDIS_EMBEDDED)s
autorestart=true
user=appuser
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=10

; ---------------------------------------------------------------------------
; Embedded PostgreSQL — skipped when POSTGRES_EMBEDDED=false
; ---------------------------------------------------------------------------
[program:postgres]
command=/usr/lib/postgresql/17/bin/postgres
    -D /data/postgres
    -c listen_addresses=127.0.0.1
    -c max_connections=50
    -c shared_buffers=128MB
    -c work_mem=4MB
    -c log_min_messages=warning
autostart=%(ENV_POSTGRES_EMBEDDED)s
autorestart=true
user=postgres
stopsignal=INT
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=20
startsecs=3

; ---------------------------------------------------------------------------
; One-shot PostgreSQL init (create user + database on first run)
; ---------------------------------------------------------------------------
[program:pg-init]
command=/app/docker/pg-init.sh
autostart=%(ENV_POSTGRES_EMBEDDED)s
autorestart=false
startsecs=0
exitcodes=0
user=root
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=25

; ---------------------------------------------------------------------------
; Caddy TLS reverse proxy — skipped when CADDY_ENABLED=false
; Terminates TLS on :443, proxies to app on :8000
; ---------------------------------------------------------------------------
[program:caddy]
command=caddy run --config /tmp/Caddyfile
autostart=%(ENV_CADDY_ENABLED)s
autorestart=true
user=root
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=90

; ---------------------------------------------------------------------------
; Agnostic QA Platform (Chainlit + FastAPI)
; Always listens on 127.0.0.1:8000 when Caddy is active, 0.0.0.0:8000 otherwise
; ---------------------------------------------------------------------------
[program:agnostic]
command=chainlit run webgui/app.py --host 0.0.0.0 --port 8000
directory=/app
user=appuser
environment=HOME="/home/appuser"
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=100
startsecs=5
